# suro-oaas/suro-oaas-devops/environments/sl-production/meteor-suro.conf
# meteorjs - meteorjs job file
# reference: http://stackoverflow.com/questions/15657071/how-to-run-meteor-on-startup-on-ubuntu-server

description "Meteor - Surgical Unit Resource Optimisation"
author "Christian Vecchiola"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Essentially lets upstart know the process will detach itself to the background
expect fork

# Start the process
# --settings and --port is optional
exec su - ubuntu -c "cd /opt/ibm/meteor && meteor --port 3001"